# backup.service
[Unit]
Description=Launches Daily Backup for Remote Shares

[Service]
Type=simple
ExecStart=/usr/bin/sh /opt/backup/backup.sh
